# api.ini
[program:api]
command=gunicorn -c gunicorn.conf wsgi:app
directory=CodeNutAPI
stdout_logfile=log/api.out
user=root
process_name=root
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
user=root