# celery.ini
[program:judge]
command=celery -A dispatcher worker -Q judge -l INFO
directory=CodeNutJudge
user=root
numprocs=1
stdout_logfile=log/worker.log
stderr_logfile=log/worker.log
autostart=true
autorestart=true
startsecs=1
stopwaitsecs = 6
killasgroup=true